# deploy/florence-api.service
[Unit]
Description=AI Nurse Florence API
After=network.target

[Service]
Type=simple
User=florence
Group=florence
WorkingDirectory=/opt/florence
EnvironmentFile=/etc/default/florence
ExecStart=/opt/venv/bin/gunicorn -k uvicorn.workers.UvicornWorker --workers 4 --bind 0.0.0.0:8000 app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
